# units/etc/sv/rtpengine-recording/run
#!/bin/bash
set -m

exec /usr/sbin/rtpengine-recording --config-file=/etc/rtpengine/rtpengine-recording.conf